# scripts/usb-automount.rules
# Udev rules for automatic USB drive mounting
# Uses systemd to run mount script with full privileges

# When a USB partition is added, trigger systemd service
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_USAGE}=="filesystem", ENV{ID_BUS}=="usb", TAG+="systemd", ENV{SYSTEMD_WANTS}="usb-mount@%k.service"

# When a USB partition is removed, run cleanup script
# This cleans up ALL USB mounts since device names can change (sda1→sdb1)
ACTION=="remove", SUBSYSTEM=="block", ENV{ID_BUS}=="usb", RUN+="/bin/systemctl stop slitcam-recorder.service", RUN+="/usr/local/bin/usb-mount.sh unmount"
